# pkg/rpm/postinstall
if [ $1 -eq 1 ] ; then
        # Initial installation
        /usr/bin/systemctl preset serviced >/dev/null 2>&1 || :
        /usr/bin/systemctl enable serviced >/dev/null 2>&1 || :
fi

mkdir -p /var/log/serviced
chgrp root /var/log/serviced
chmod 1750 /var/log/serviced

#
# NOTE: changing ownership/permissions here has the side-effect of causing
#       "rpm -V serviced" to complain, but we could not get fpm to assign
#       the ownership/permissions at build time.
#
chgrp serviced /etc/default/serviced
chmod 640 /etc/default/serviced

chgrp serviced /opt/serviced
chmod 750 /opt/serviced

LIBDEVMAPPER="$(ldconfig -p | grep libdevmapper.so.1.02 | awk {'print $4'})"
if [ -z "${LIBDEVMAPPER}" ]; then
    exit 0
fi
if ! (python -c "import sys; sys.exit(0 if \"${LIBDEVMAPPER}\".endswith('.1') else 1)"); then
    [ -f "${LIBDEVMAPPER}.1" -o -h "${LIBDEVMAPPER}.1" ] || ln -sf ${LIBDEVMAPPER} ${LIBDEVMAPPER}.1
    ldconfig
fi
echo vm.max_map_count=262144 >> /etc/sysctl.conf
sysctl --system
