# pkg/cron.daily
#!/bin/sh

/usr/sbin/logrotate /opt/serviced/etc/logrotate.conf
/opt/serviced/bin/serviced-container-cleanup
/opt/serviced/bin/serviced-btrfs-maintenance
