# pkg/cron.daily
#!/bin/sh

/usr/sbin/logrotate /opt/serviced/etc/logrotate.conf
